# install/remoteAppStandAlone.service
[Unit]
Description=remoteAppStandAlone

[Service]
Type=simple
User=volumio
Group=volumio
ExecStart=/usr/local/bin/remoteAppStandAlone start
ExecStop=/usr/local/bin/remoteAppStandAlone stop


[Install]
WantedBy=multi-user.target